# src/systemd/lima_ddns.service
[Unit]
Description = Updates Lima-City DDNS records

[Service]
# DynamicUser = yes
ExecStart = /usr/bin/lima_ddns
Type = oneshot
